# s6-overlay/s6-rc.d/gsd-web/run
#!/bin/sh
cd /workspace
exec with-contenv s6-setuidgid claude env HOME=/home/claude /usr/local/bin/gsd-web.sh
