# s6-overlay/s6-rc.d/gateway/run
#!/bin/sh
exec node /usr/local/share/holyclaude/gateway.js
